# server/flask/application/home-vod-server.conf
description "uWSGI instance to serve home-vod-server"

start on runlevel [2345]
stop on runlevel [!2345]

setgid www-data

script
    cd $HOME/server
    uwsgi --ini home-vod-server.ini
 end-script